# build/foo-1.2.3-4_amd64/lib/systemd/system/foo.service
[Unit]
Description=FOO Service
After=network.target

[Service]
ExecStart=/opt/foo/bin/foo.sh
WorkingDirectory=/var/foo
KillSignal=SIGINT
Restart=on-failure
User=foo
StandardOutput=null
StandardError=journal
SyslogIdentifier=foo

[Install]
WantedBy=multi-user.target
