# build/foo-1.2.3-4_amd64/DEBIAN/prerm
#!/bin/bash
set -e

# Stop and disable the service
systemctl stop foo.service
systemctl disable foo.service

